# systemd/riskline.timer
[Unit]
Description=Run Riskline every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true
Unit=riskline.service

[Install]
WantedBy=timers.target
